# .config/systemd/user/polybar.service
[Unit]
Description=A fast and easy-to-use tool status bar
Documentation=man:polybar(1)
Wants=autorandr.service
After=autorandr.service

[Service]
Type=simple
ExecStart=%h/.config/polybar/launch.sh %U %h
Restart=always

[Install]
WantedBy=default.target
