# .config/systemd/user/sxhkd.service
[Unit]
Description=Simple X Hotkey Daemon
Documentation=man:sxhkd(1)
PartOf=graphical-session.target
Wants=graphical-session.target sxhkd-fifo.service 
After=graphical-session.target sxhkd-fifo.service

[Service]
Environment=DISPLAY=:0 SHELL=bash
ExecStart=/usr/bin/sxhkd -s /tmp/sxhkd.fifo
ExecReload=/usr/bin/kill -SIGUSR1 $MAINPID
Restart=on-failure

[Install]
WantedBy=default.target
