# .config/systemd/user/sxhkd.service
[Unit]
Description=Simple X Hotkey Daemon
Documentation=man:sxhkd(1)
# BindsTo=display-manager.service
# After=display-manager.service

[Service]
Environment=DISPLAY=:0
ExecStartPre=rm -f /tmp/sxhkd.fifo && mkfifo /tmp/sxhkd.fifo
ExecStart=/usr/bin/sxhkd -s /tmp/sxhkd.fifo
ExecReload=/usr/bin/kill -SIGUSR1 $MAINPID
Restart=on-failure

[Install]
WantedBy=default.target
